# main/alpine-baselayout/inittab
# /etc/inittab

::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

# calls a script named "alpine-linux-initialize.sh" that could be either
# inside a USB stick or a CDROM. The idea is to allow the user of live
# Alpine Linux to quickly modify the .iso and add a script than can 
# perform whatever initialization he requires whenever the live system
# boots.
::sysinit:if [ -e "/media/usb/alpine-linux-initialize.sh" ]; then "/media/usb/alpine-linux-initialize.sh else "/media/cdrom/alpine-linux-initialize.sh" fi

# Set up a couple of getty's
tty1::respawn:/sbin/getty 38400 tty1
tty2::respawn:/sbin/getty 38400 tty2
tty3::respawn:/sbin/getty 38400 tty3
tty4::respawn:/sbin/getty 38400 tty4
tty5::respawn:/sbin/getty 38400 tty5
tty6::respawn:/sbin/getty 38400 tty6

# Put a getty on the serial port
#ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/sbin/openrc shutdown

